# vagrant/flower_upstart.conf
description "Celery flower"
start on vagrant-mounted
stop on runlevel [016]
exec /home/vagrant/.virtualenvs/edx-ora2/bin/celery flower
respawn
respawn limit 15 5
